# vps_setup/email-webhook.service
[Unit]
Description=Email Webhook Bridge (Maddy to Railway)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/email_webhook.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
